# packages/addons/driver/tbs-linux-cc/source/system.d/driver.dvb.tbs-linux-cc.service
[Unit]
Description=tbs-linux-cc driver service
After=network.target

[Service]
Type=oneshot
ExecStart=-/bin/sh -c "exec sh /storage/.kodi/addons/driver.dvb.tbs-linux-cc/bin/tbs-linux-cc.start"
ExecStop=-/bin/sh -c "exec sh /storage/.kodi/addons/driver.dvb.tbs-linux-cc/bin/tbs-linux-cc.stop"
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
